# assets/nomad.ini
[program:nomad]
command=/usr/local/bin/nomad agent -config=/etc/nomad.d
autostart=true
autorestart=true
stderr_logfile=/var/log/nomad.err.log
stdout_logfile=/var/log/nomad.out.log

[program:consul]
command=/usr/local/bin/consul agent -config-file=/etc/consul.d/server.hcl
autostart=true
autorestart=true
stderr_logfile=/var/log/consul.err.log
stdout_logfile=/var/log/consul.out.log

[program:docker]
command=/usr/local/bin/dockerd-entrypoint.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/docker.err.log
stdout_logfile=/var/log/docker.out.log